# polybar/.config/systemd/user/polybar.service
[Unit]
Description=Polybar status bar

[Service]
Type=forking
ExecStart=/home/derek/.config/polybar/launch.sh
