# sxhkd/.config/systemd/user/sxhkd@bspwm.service
[Unit]
Description=A simple X hotkey daemon

[Service]
ExecStart=sxhkd /home/derek/.config/bspwm/bspwm.sxhkdrc
