# udiskie/.config/autostart/udiskie.desktop
[Desktop Entry]
Type=Application
NoDisplay=false
Name=udiskie
Categories=Utility;
TryExec=udiskie
Exec=udiskie --smart-tray
