# udiskie/.config/autostart/udiskie.desktop
[Desktop Entry]
Name=Udiskie
Exec=udiskie --appindicator
Type=Application
