# X11/.config/systemd/user/keyboard.service
[Unit]
Description=Keyboard Configuration

[Service]
ExecStart=/home/derek/.local/bin/setkb
