# X11/.config/systemd/user/numlockx.service
[Unit]
Description=Numlock activation

[Service]
ExecStart=numlockx
