# X11/.config/systemd/user/unclutter.service
[Unit]
Description=Hides idle pointer

[Service]
ExecStart=unclutter
