# X11/.config/systemd/user/wallpaper.service
[Unit]
Description=Background Image

[Service]
ExecStart=/home/derek/.local/bin/setbg
