# picom/.config/systemd/user/picom.service
[Unit]
Description=picom X11 compositor

[Service]
ExecStart=picom
